# shortcuts/.config/autostart/steam.desktop
[Desktop Entry]
Type=Application
Name=Steam
Exec=steam -silent
Icon=steam
Terminal=false